# vars/logstash/templates/logstash-upstart.conf.j2
description "Logstash container"
author "Mizeraj"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
  /usr/bin/docker start -a {{ logstash_container_name }}
end script
